# gh-trending-tracker.service.example
[Unit]
Description=GitHub Trending Tracker Service
Documentation=https://github.com/yourusername/gh-trending-tracker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# User and group to run the service as
User=www-data
Group=www-data

# Working directory
WorkingDirectory=/opt/gh-trending-tracker

# Virtual environment and Python interpreter
Environment="PATH=/opt/gh-trending-tracker/venv/bin:/usr/local/bin:/usr/bin:/bin"

# Proxy configuration (uncomment if needed)
#Environment="HTTP_PROXY=http://proxy.example.com:8080"
#Environment="HTTPS_PROXY=http://proxy.example.com:8080"
#Environment="NO_PROXY=localhost,127.0.0.1"

ExecStart=/opt/gh-trending-tracker/venv/bin/python -m src.scheduler

# Restart policy
Restart=always
RestartSec=10

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/gh-trending-tracker/data /opt/gh-trending-tracker/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gh-trending-tracker

[Install]
WantedBy=multi-user.target
